# Temp/tempo/home/pi/.config/autostart/programme.desktop
[Desktop Entry]
Exec=lxterminal –e "sudo python /home/pi/superGSM.py"
Type=Application
Terminal=true
X-KeepTerminal=true
StartupNotify=true
