# Temp/tempo/home/pi/.config/autostart/vncDaZe.desktop
[Desktop Entry]
Exec=lxterminal -l –e "/usr/bin/tightvncserver"
Type=Application
Terminal=true
X-KeepTerminal=true
StartupNotify=true
